# services/healthchecker/healthchecker.service
[Unit]
Description=Check the health of regitered services
After=mysql.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash $script_path $logging_dir $venv_path

[Install]
WantedBy=multi-user.target